# deploy/gpa-v3-icloud-contacts-sync.timer
[Unit]
Description=Run GPA V3 iCloud contacts sync every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=gpa-v3-icloud-contacts-sync.service

[Install]
WantedBy=timers.target
