# deploy/gpa-v3-telegram-reminder.service
[Unit]
Description=GPA V3 daily Telegram reminder
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=gpa
Group=gpa
WorkingDirectory=/opt/gpa-v3/app
EnvironmentFile=/opt/gpa-v3/app/.env
ExecStart=/opt/gpa-v3/venv/bin/python -m backend.jobs.telegram_reminder
