# deploy/gpa-v3-icloud-contacts-sync.service
[Unit]
Description=GPA V3 iCloud contacts sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=gpa
Group=gpa
WorkingDirectory=/opt/gpa-v3/app
EnvironmentFile=/opt/gpa-v3/app/.env
Environment=GPA_ICLOUD_CONTACTS_FILE=/opt/gpa-v3/secrets/icloud_contacts.env
ExecStart=/opt/gpa-v3/venv/bin/python -m backend.jobs.icloud_contacts_sync
